# kryptofranccore/contrib/init/bitcoind.conf
description "KryptoFranc Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env KRYPTOFRANCD_BIN="/usr/bin/kryptofrancd"
env KRYPTOFRANCD_USER="kryptofranc"
env KRYPTOFRANCD_GROUP="kryptofranc"
env KRYPTOFRANCD_PIDDIR="/var/run/kryptofrancd"
# upstart can't handle variables constructed with other variables
env KRYPTOFRANCD_PIDFILE="/var/run/kryptofrancd/kryptofrancd.pid"
env KRYPTOFRANCD_CONFIGFILE="/etc/kryptofranc/kryptofranc.conf"
env KRYPTOFRANCD_DATADIR="/var/lib/kryptofrancd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # kryptofrancd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$KRYPTOFRANCD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run kryptofrancd."
        echo "The setting must appear in $KRYPTOFRANCD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"kryptoFranc Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$KRYPTOFRANCD_PIDDIR"
    chmod 0755 "$KRYPTOFRANCD_PIDDIR"
    chown $KRYPTOFRANCD_USER:$KRYPTOFRANCD_GROUP "$KRYPTOFRANCD_PIDDIR"
    chown $KRYPTOFRANCD_USER:$KRYPTOFRANCD_GROUP "$KRYPTOFRANCD_CONFIGFILE"
    chmod 0660 "$KRYPTOFRANCD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$KRYPTOFRANCD_PIDFILE" \
    --chuid $KRYPTOFRANCD_USER:$KRYPTOFRANCD_GROUP \
    --exec "$KRYPTOFRANCD_BIN" \
    -- \
    -pid="$KRYPTOFRANCD_PIDFILE" \
    -conf="$KRYPTOFRANCD_CONFIGFILE" \
    -datadir="$KRYPTOFRANCD_DATADIR" \
    -disablewallet \
    -daemon

